# install/etc/s6/services/10-postal/run
#!/usr/bin/with-contenv bash

source /assets/functions/*
PROCESS_NAME="postal"
check_initialized
liftoff

print_info "Starting postal"
exec /opt/postal/bin/postal run
